# ch-backup-weekly.service
[Unit]
Description=Weekly ClickHouse backup of all user tables
After=clickhouse-server.service
Requires=clickhouse-server.service

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/sbin/ch-backup-all --local
StandardOutput=append:/var/log/wsprdaemon/ch-backup-weekly.log
StandardError=append:/var/log/wsprdaemon/ch-backup-weekly.log
TimeoutStartSec=6h
