# pskreporter-forward-elector.service
# /etc/systemd/system/pskreporter-forward-elector.service
#
# gw1-side elector for the PSKReporter forwarder (Phase 2 PR 2).
# Probes wd10/wd20/wd30 ClickHouse ports and writes the preferred
# healthy leader's short hostname to /var/www/html/pskreporter-leader.txt
# (served by nginx at http://gw1.wsprdaemon.org/pskreporter-leader.txt).
[Unit]
Description=PSKReporter forwarder elector
Documentation=https://github.com/rrobinett/wsprdaemon-server
After=network-online.target nginx.service
Wants=network-online.target

[Service]
Type=simple
User=wsprdaemon
Group=www-data
EnvironmentFile=-/etc/default/pskreporter-forward-elector
ExecStart=/opt/wsprdaemon-server/venv/bin/python \
    /opt/wsprdaemon-server/pskreporter_forward_elector.py
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pskreporter-forward-elector

[Install]
WantedBy=multi-user.target
